# supervisord.conf
[supervisord]
nodaemon=true

[program:kestrel]
directory=/VNextApp  
command=k kestrel  
